# sources/rtorrent.service
[Unit]
Description=Flood Web UI
After=network.target

[Service]
WorkingDirectory=~
ExecStart=/usr/bin/flood
User=flood

[Install]
WantedBy=multi-user.target